# etc/systemd/system/android-mount.service.d/30-modules.conf
[Service]
ExecStartPost=-/usr/bin/setup-modules.sh
