# etc/systemd/system/android-mount.service.d/40-backlight.conf
[Service]
ExecStartPost=-/usr/bin/rm -rf /var/lib/systemd/backlight
